# dist/lib/systemd/system/homeassistant.service
[Unit]
Description=Home Assistant
After=network-online.target

[Service]
Type=simple
User=%i
ExecStart=/srv/homeassistant/bin/hass -c "/var/local/homeassistant/.homeassistant/"

[Install]
WantedBy=multi-user.target
